# root/etc/s6-overlay/s6-rc.d/svc-dbus/run
#!/usr/bin/with-contenv bash

# Run dbus
mkdir /run/dbus
dbus-daemon \
  --system \
  --nofork \
  --nosyslog 
