# root/etc/s6-overlay/s6-rc.d/svc-selkies/run
#!/usr/bin/with-contenv bash

# Clear the cache registry
rm -rf "${HOME}/.cache/gstreamer-1.0"

# Set default resolution
s6-setuidgid abc selkies-gstreamer-resize "${DISPLAY_SIZEW}x${DISPLAY_SIZEH}"


# Setup dev mode if defined
if [ ! -z ${DEV_MODE+x} ]; then
  # Frontend setup
  apt-get update
  apt-get install -y \
    nodejs
  npm install -g nodemon
  cd $HOME/src/addons/gst-web-core
  s6-setuidgid abc npm install
  s6-setuidgid abc npm run serve &
  # Run backend
  cd $HOME/src/src/selkies_gstreamer
  s6-setuidgid abc \
    nodemon --exec \
      "python3" __main__.py \
      --addr="localhost" \
      --port="8081" \
      --enable_basic_auth="false" \
      --enable_metrics_http="true" \
      --metrics_http_port="9081"
fi


# Start the Selkies-GStreamer WebRTC HTML5 remote desktop application
exec s6-setuidgid abc \
  selkies-gstreamer \
    --addr="localhost" \
    --port="8081" \
    --enable_basic_auth="false" \
    --enable_metrics_http="true" \
    --metrics_http_port="9081" 
